# docker/config/supervisor/conf.d/avahi-daemon.conf
[program:avahi-daemon]
command=/usr/sbin/avahi-daemon --no-chroot
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0