# docker/config/supervisor/conf.d/smbd.conf
[program:smbd]
command=/usr/sbin/smbd --no-process-group -F -d 1 -S -s /etc/samba/smb.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0